# archiso/airootfs/etc/systemd/system/systemd-firstboot.service.d/override.conf
# Disable systemd-firstboot on Live ISO
# VV OS has its own installer that handles timezone/locale setup

[Service]
ExecStart=
ExecStart=/usr/bin/true
